# root/etc/s6-overlay/s6-rc.d/digrouz-base-adduser/run
#!/command/with-contenv sh

UID=${UID:-12345}
GID=${GID:-12345}

groupmod -o -g "${GID}" appuser
usermod -o -u "${UID}" appuser